# sysutils/powernowd/services/powernowd.service
[Unit]
Description=PowerNow! Daemon

[Service]
Type=simple
ExecStart=/usr/bin/powernowd -d -q

[Install]
WantedBy=multi-user.target
